# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres-err.log
priority=10

[program:redis]
command=redis-server --maxmemory 256mb --maxmemory-policy allkeys-lru --daemonize no
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-err.log
priority=20

[program:ollama]
command=ollama serve
autostart=true
autorestart=true
environment=OLLAMA_HOST="0.0.0.0:11434"
stdout_logfile=/var/log/supervisor/ollama.log
stderr_logfile=/var/log/supervisor/ollama-err.log
priority=30

[program:oak-api-proxy]
command=uvicorn main:app --host 0.0.0.0 --port 9000
directory=/opt/oak/proxy
environment=OLLAMA_BASE_URL="http://localhost:11434",OLLAMA_DEFAULT_MODEL="qwen3-coder",REDIS_URL="redis://localhost:6379",ROUTING_STRATEGY="passthrough"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/proxy.log
stderr_logfile=/var/log/supervisor/proxy-err.log
priority=40
startsecs=5

[program:oak-api]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/opt/oak
environment=DATABASE_URL="postgresql://oak:oak@localhost:5432/oak",REDIS_URL="redis://localhost:6379",PYTHONPATH="/opt/oak"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api-err.log
priority=50
startsecs=10

[program:oak-ui]
command=streamlit run app.py --server.port 8501 --server.headless true --server.address 0.0.0.0
directory=/opt/oak/ui
environment=OAK_API_URL="http://localhost:8000"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ui.log
stderr_logfile=/var/log/supervisor/ui-err.log
priority=60
startsecs=5
